# day-30-private-cosmetics-service-kotlin/deploy/ollama.service.d/override.conf
[Service]
# Never expose Ollama directly. Only the Kotlin service may call it over loopback.
Environment="OLLAMA_HOST=127.0.0.1:11434"
Environment="OLLAMA_MAX_LOADED_MODELS=1"
Environment="OLLAMA_NUM_PARALLEL=1"
Environment="OLLAMA_MAX_QUEUE=2"
Environment="OLLAMA_CONTEXT_LENGTH=8192"
Environment="OLLAMA_NOHISTORY=1"
Restart=on-failure
RestartSec=5s
